# packetchess@.service
[Unit]
Description=PacketChess Server

[Service]
ExecStart=/opt/venv/bin/python /home/bpquser/packetchess.py
StandardInput=socket
StandardOutput=socket
User=bpquser

[Install]
WantedBy=multi-user.target
